# deploy/systemd/mywebapp.socket
[Unit]
Description=My Web App Socket

[Socket]
ListenStream=127.0.0.1:8080
Accept=no

[Install]
WantedBy=sockets.target